# root/etc/services.d/sonarr/run
#!/usr/bin/with-contenv bash

umask 0002

XDG_CONFIG_HOME="/config/xdg" exec s6-setuidgid abc mono /opt/NzbDrone/NzbDrone.exe -nobrowser -data=/config
