# root/etc/services.d/sonarr/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}

umask "$UMASK_SET"

cd /opt/NzbDrone || exit

exec \
	s6-setuidgid abc mono --debug NzbDrone.exe \
	-nobrowser -data=/config
