# root/etc/services.d/sonarr/run
#!/usr/bin/with-contenv bash

cd /app/sonarr/bin || exit

exec \
    s6-setuidgid abc mono --debug Sonarr.exe \
    -nobrowser -data=/config
