# root/etc/services.d/sonarr/run
#!/usr/bin/with-contenv bash

umask 022

cd /opt/sonarr || exit

exec \
	s6-setuidgid abc mono --debug NzbDrone.exe \
	-nobrowser -data=/config
