# root/etc/services.d/sonarr/run
#!/usr/bin/with-contenv bash

umask 0002

cd /opt/NzbDrone

exec \
	s6-setuidgid abc mono NzbDrone.exe \
	-nobrowser -data=/config
