# root-armhf/etc/services.d/armless/run
#!/usr/bin/with-contenv bash

if [[ -n $ISCI ]]; then
  sleep 180
fi

exec \
    s6-svscanctl -t /run/service/
