# root-armhf/etc/services.d/armless/run
#!/usr/bin/with-contenv bash

exec \
    s6-svscanctl -t /run/service/
