# services/sonarr/run
#!/bin/bash

exec /sbin/setuser abc mono /opt/NzbDrone/NzbDrone.exe -nobrowswer -data=/config
