# services/sonarr/run
#!/bin/bash
umask 0002
exec /sbin/setuser abc mono /opt/NzbDrone/NzbDrone.exe -nobrowswer -data=/config
